# .docker/s6-rc.d/php/run
#!/bin/sh -e

/usr/local/sbin/php-fpm --nodaemonize
